# watcher/openclaw-watcher.service
[Unit]
Description=OpenClaw Agent Watcher
Documentation=https://github.com/openclaw/openclaw
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/.openclaw/workspace/watcher
ExecStart=/usr/bin/node /root/.openclaw/workspace/watcher/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
